# pyodin/sys/docs/supervisor.ini
[program:odinsys]
command=/home/vassilux/Projects/nodeProjects/odin/odinsys/run_odin_sysmon.py
environment=PATH="/home/vassilux/Projects/virtualEnvs/odin/bin"
user=vassilux
autostart=true
autorestart=true
